# wifi-watchdog.timer
[Unit]
Description=WiFi Connectivity Watchdog Timer
After=network-online.target

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min
Unit=wifi-watchdog.service
Persistent=true

[Install]
WantedBy=timers.target
